# debian/nfblock.cron.daily
#!/bin/sh

ENABLED=0
BLOCKLIST_URL=

LIST_PATH=/var/lib/nfblock

test -f /etc/default/nfblock && . /etc/default/nfblock

if [ "$ENABLED" = "0" ]; then
    exit 0
fi


if [ -n "$BLOCKLIST_URL" ] ; then
    for N in $BLOCKLIST_URL ; do
        wget -q -N -P "$LIST_PATH" "$N"
    done

    invoke-rc.d nfblock reload
fi
